# etc/exabgpmon.conf
[program:exabgp]
command=exabgp etc/conf.ini
autostart=false                ; start at supervisord start (default: true)
autorestart=unexpected        ; whether/when to restart (default: unexpected)
